# debian/ceilometer-agent-compute.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
        update-rc.d ceilometer-agent-compute remove
fi

#DEBHELPER#
